# raspberry-pi-spotify-player/spotify-player.service
[Unit]
Description=Spotify Continuous Player
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/spotify-player
ExecStart=/usr/bin/python3 /home/pi/spotify-player/spotify_player.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
